# r-base.conf
[program:r-base]
command=/bin/bash -c "R"
redirect_stderr=true
