# r-base.conf
[program:r-base]
command=/usr/bin/R
redirect_stderr=true
